# src/Unoserver/PoolManager/unoserver@.service
[Unit]
Description=Unoserver instance on port %i
After=network.target
Documentation=https://github.com/unoconv/unoserver

[Service]
Type=simple
EnvironmentFile=-/etc/unoserver/env.conf
ExecStart=/usr/bin/unoserver --port %i ${UNOSERVER_EXTRA_OPTS}
Restart=always
RestartSec=5
User=${UNOSERVER_USER:-nobody}
Group=${UNOSERVER_GROUP:-nogroup}
# Limitar recursos
MemoryMax=2G
CPUQuota=50%
# Seguridad
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target